# scratch.crontab
#55 18 * * * rsync -Parv /nird/datapeak/NS9036K/promec/promec/Qexactive/Raw/ /nird/datalake/NS9036K/orbitrap/
#55 09 * * * rsync -Parv /nird/datapeak/NS9036K/promec/promec/TIMSTOF/Raw/ /nird/datalake/NS9036K/Raw/
05 00 * * *  find /projects/NS9036K/NORSTORE_OSL_DISK/NS9036K/promec/promec -mtime 0 -type f -print0 | tar -czvf /projects/NS9036K/NORSTORE_OSL_DISK/NS9036K/promec/backup.$(date +"\%m_\%d_\%Y").tar.gz --null -T - > /projects/NS9036K/NORSTORE_OSL_DISK/NS9036K/promec/backup.$(date +"\%m_\%d_\%Y").log 2>/projects/NS9036K/NORSTORE_OSL_DISK/NS9036K/promec/backup.2.$(date +"\%m_\%d_\%Y").log 1>/projects/NS9036K/NORSTORE_OSL_DISK/NS9036K/promec/backup.1.$(date +"\%m_\%d_\%Y").log
20 15 * * * rclone-v1.67.0-linux-amd64/rclone sync --drive-pacer-burst 50 --drive-pacer-min-sleep 1s  -Pv DL/Raw/ remote:TIMSTOF/ 0>driveTIMSTOF.0.log 1>driveTIMSTOF.1.log 2>driveTIMSTOF.2.log > driveTIMSTOF.log
19 15 * * * killall rclone
55 23 * * * rclone-v1.67.0-linux-amd64/rclone sync --drive-stop-on-upload-limit --tpslimit 10   -Pv DL/orbitrap/ remote:orbitrap/ 0>driveorbi.0.log 1>driveorbi.1.log 2>driveorbi.2.log > driveorbi.log


