# scratch.crontab
#58 * * * * flock -w 0 $HOME/gdrive.lock /home/animeshs/gopath/bin/drive push -ignore-name-clashes -verbose --no-prompt $HOME/promec 0>>drive.0.log 1>>drive.1.log 2>>drive.2.log >> drive.log
00 * * * * flock -w 0 $HOME/cron.lock rsync -Parv --min-size=1  --exclude=.gd --exclude=.tmp.drivedownload  /home/animeshs/promec/ ash022@login.nird.sigma2.no:promec/NORSTORE_OSL_DISK/NS9036K/promec/  2>cronsync.2.log 1>cronsync.1.log > cronsync.log
30 * * * *  ssh -n -f  ash022@login.nird.sigma2.no "nohup bash ./scripts/dataRaw.sh  >> nohupDataRaw.out &"
51 18 * * *  ssh -n -f  ash022@login.nird.sigma2.no "find /projects/NS9036K/NORSTORE_OSL_DISK/NS9036K/promec/promec -mtime 0 -type f -print0 | tar -czvf /projects/NS9036K/NORSTORE_OSL_DISK/NS9036K/promec/backup.$(date +"\%m_\%d_\%Y").tar.gz --null -T - > /projects/NS9036K/NORSTORE_OSL_DISK/NS9036K/promec/backup.$(date +"\%m_\%d_\%Y").log 2>/projects/NS9036K/NORSTORE_OSL_DISK/NS9036K/promec/backup.2.$(date +"\%m_\%d_\%Y").log 1>/projects/NS9036K/NORSTORE_OSL_DISK/NS9036K/promec/backup.1.$(date +"\%m_\%d_\%Y").log"
#10 * * * * flock -w 0 $HOME/cron.lock  /home/animeshs/gopath/bin/drive push --no-prompt /home/animeshs/Documents 2>>cronsync.2.log 1>>cronsync.1.log >> cronsync.log
