# scratch.crontab
00 * * * * flock -w 0 $HOME/cron.lock rsync -Parv --min-size=1  --exclude=.gd --exclude=promec/.tmp.driveupload  /home/animeshs/promec/ ash022@login-trd.nird.sigma2.no:promec/NORSTORE_OSL_DISK/NS9036K/promec/  2>cronsync.2.log 1>cronsync.1.log > cronsync.log
30 * * * *  ssh -n -f  ash022@login-trd.nird.sigma2.no "nohup bash ./scripts/dataRaw.sh  >> nohupDataRaw.out &"
51 18 * * *  ssh -n -f  ash022@login-trd.nird.sigma2.no "find /projects/NS9036K/NORSTORE_OSL_DISK/NS9036K/promec/promec -mtime 0 -type f -print0 | tar -czvf /projects/NS9036K/NORSTORE_OSL_DISK/NS9036K/promec/backup.$(date +"\%m_\%d_\%Y").tar.gz --null -T - > /projects/NS9036K/NORSTORE_OSL_DISK/NS9036K/promec/backup.$(date +"\%m_\%d_\%Y").log 2>/projects/NS9036K/NORSTORE_OSL_DISK/NS9036K/promec/backup.2.$(date +"\%m_\%d_\%Y").log 1>/projects/NS9036K/NORSTORE_OSL_DISK/NS9036K/promec/backup.1.$(date +"\%m_\%d_\%Y").log"
30 9 * * 5  ssh -n -f  ash022@login0-trd.nird.sigma2.no  "killall drive"
35 9 * * 5  ssh -n -f  ash022@login0-trd.nird.sigma2.no  "nohup /nird/home/ash022/gopath/gopath/bin/drive push  --no-prompt --ignore-name-clashes -verbose -force /projects/NS9036K/NORSTORE_OSL_DISK/NS9036K/promec/promec/HF/  > /projects/NS9036K/NORSTORE_OSL_DISK/NS9036K/promec/HF.$(date +"\%m_\%d_\%Y").log 2>/projects/NS9036K/NORSTORE_OSL_DISK/NS9036K/promec/HF.2.$(date +"\%m_\%d_\%Y").log 1>/projects/NS9036K/NORSTORE_OSL_DISK/NS9036K/promec/HF.1.$(date +"\%m_\%d_\%Y").log"
40 9 * * 5  ssh -n -f  ash022@login0-trd.nird.sigma2.no  "nohup /nird/home/ash022/gopath/gopath/bin/drive push  --no-prompt --ignore-name-clashes -verbose -force /projects/NS9036K/NORSTORE_OSL_DISK/NS9036K/promec/promec/TIMSTOF/  > /projects/NS9036K/NORSTORE_OSL_DISK/NS9036K/promec/TIMSTOF.$(date +"\%m_\%d_\%Y").log 2>/projects/NS9036K/NORSTORE_OSL_DISK/NS9036K/promec/TIMSTOF.2.$(date +"\%m_\%d_\%Y").log 1>/projects/NS9036K/NORSTORE_OSL_DISK/NS9036K/promec/TIMSTOF.1.$(date +"\%m_\%d_\%Y").log &"
45 9 * * 5  ssh -n -f  ash022@login0-trd.nird.sigma2.no  "nohup /nird/home/ash022/gopath/gopath/bin/drive push  --no-prompt --ignore-name-clashes -verbose -force /projects/NS9036K/NORSTORE_OSL_DISK/NS9036K/promec/promec/Elite/  > /projects/NS9036K/NORSTORE_OSL_DISK/NS9036K/promec/Elite.$(date +"\%m_\%d_\%Y").log 2>/projects/NS9036K/NORSTORE_OSL_DISK/NS9036K/promec/Elite.2.$(date +"\%m_\%d_\%Y").log 1>/projects/NS9036K/NORSTORE_OSL_DISK/NS9036K/promec/Elite.1.$(date +"\%m_\%d_\%Y").log"
50 9 * * 5  ssh -n -f  ash022@login0-trd.nird.sigma2.no  "nohup /nird/home/ash022/gopath/gopath/bin/drive push  --no-prompt --ignore-name-clashes -verbose -force /projects/NS9036K/NORSTORE_OSL_DISK/NS9036K/promec/promec/Qexactive/  > /projects/NS9036K/NORSTORE_OSL_DISK/NS9036K/promec/Qexactive.$(date +"\%m_\%d_\%Y").log 2>/projects/NS9036K/NORSTORE_OSL_DISK/NS9036K/promec/Qexactive.2.$(date +"\%m_\%d_\%Y").log 1>/projects/NS9036K/NORSTORE_OSL_DISK/NS9036K/promec/Qexactive.1.$(date +"\%m_\%d_\%Y").log"
55 9 * * 5  ssh -n -f  ash022@login0-trd.nird.sigma2.no  "nohup /nird/home/ash022/gopath/gopath/bin/drive push  --no-prompt --ignore-name-clashes -verbose -force /projects/NS9036K/NORSTORE_OSL_DISK/NS9036K/promec/promec/USERS/  > /projects/NS9036K/NORSTORE_OSL_DISK/NS9036K/promec/USERS.$(date +"\%m_\%d_\%Y").log 2>/projects/NS9036K/NORSTORE_OSL_DISK/NS9036K/promec/USERS.2.$(date +"\%m_\%d_\%Y").log 1>/projects/NS9036K/NORSTORE_OSL_DISK/NS9036K/promec/USERS.1.$(date +"\%m_\%d_\%Y").log"
#10 * * * * flock -w 0 $HOME/cron.lock  /home/animeshs/gopath/bin/drive push --no-prompt /home/animeshs/Documents 2>>cronsync.2.log 1>>cronsync.1.log >> cronsync.log
00 07 * * * python /home/animeshs/scripts/unlock.py 2> unlock.txt
00 19 * * * python /home/animeshs/scripts/lock.py 2> lock.txt
