# scratch.crontab
33 22 * * * $HOME/scripts/dataRaw.sh 0> dataRaw.0.log 1>dataRaw.1.log 2>dataRaw.2.log > dataRaw.log

